# etc/systemd/user/alacritty-theme-switcher-restart.path
[Unit]
Description="restart alacritty-theme-switcher if source changes"

[Path]
PathModified=/home/neil/src/dotfiles/lib/alacritty-theme-switcher.sh

[Install]
WantedBy=default.target
