# etc/systemd/user/alacritty-theme-switcher-restart.service
[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart --user alacritty-theme-switcher.service
