# etc/systemd/user/wallpaper.service
[Service]
Type=oneshot
ExecStart=/home/neil/src/dotfiles/lib/rotate-background.sh
