# etc/systemd/user/wallpaper.timer
[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=default.target
