# home/dot_config/systemd/user/wallpaper.service.tmpl
[Service]
Type=oneshot
ExecStart=%h/.local/lib/rotate-background.sh
